# prototypes/proto-hub.service
[Unit]
Description=Prototype Hub Server
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/.openclaw/workspace/prototypes
ExecStart=/usr/bin/node /home/ubuntu/.openclaw/workspace/prototypes/server.mjs
Restart=always
RestartSec=2
Environment=PROTO_PORT=8095
Environment=PROTO_USER=synder
Environment=PROTO_PASS=ux-review-2026

[Install]
WantedBy=multi-user.target
